# roles/docker-storage/templates/var-lib-docker.mount
[Unit]
Description=Mount /var/lib/docker
[Mount]
What={{device}}1
Where=/var/lib/docker
Type={{fs_type}}
